# board/mbay/rpi3/80autousb.rules
KERNEL=="sd*", SUBSYSTEM=="block", RUN+="/socketpinger.sh $env{ACTION} %k"